# updates/arch-updates.service
[Unit]
Description=download package updates

[Service]
Type=simple
ExecStart=/usr/bin/yay -Syuw --noconfirm
